# services/chatbot/run
#!/bin/bash
export USER=app

su - app -c 'cd /src ; npm start'
