# services/chatbot/run
#!/bin/bash
export USER=app

su - app -c 'cd /src; /src/node_modules/.bin/sequelize db:migrate --env production; npm start'
