# services/chatbot/run
#!/bin/bash
export USER=app

su - app -c 'sequelize db:migrate --env production ; cd /src ; npm start'
