# services/chatbot/run
#!/bin/bash
export USER=app

su - app -c 'cd /src; sequelize db:migrate --env production; npm start'
